# project-spec/meta-user/recipes-apps/bootscript/files/bootscript.service
[Unit]
Description=bootscript
  
[Service]
ExecStart=/usr/bin/bootscript
StandardOutput=journal+console
  
[Install]
WantedBy=multi-user.target
